# ndhcpd.service
[Unit]
Description=DHCP Server Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=ndhcpd
Group=ndhcpd
Environment="DHCP_CONFIG=/etc/ndhcpd/config.yaml"
Environment="RUST_LOG=info"
ExecStart=/usr/local/bin/ndhcpd
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ndhcpd

# Capabilities needed for DHCP (port 67)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
